# tmp/droid-config.files
%defattr(-,root,root,-)
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/ofono/noplugin.conf
/var/lib/environment/nemo/60-multimedia.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/binder.conf
/etc/acdbdata
/etc/sensorfw/primaryuse.conf
/etc/modprobe.d/blacklist-evbug.conf
/etc/selinux/config
/etc/obexd/pbap.conf
/etc/obexd/root
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/opp.conf
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/syncevolution
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/mas
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/bin/logcat
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/99-qtscenegraph.conf
/etc/usb-moded/10-setup-configfs.ini
/etc/usb-moded/80-use-ask-mode.ini
/etc/mce/60-proximity-sensor-ondemand.conf
/etc/ofono/ril_subscription.conf
/etc/acdbdata
/etc/gst-droid/gstdroidcodec.conf
/usr/libexec/droid-hybris/system/etc/init/vndservicemanager.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/libexec/droid-hybris/system/etc/init/hwservicemanager.rc
/usr/libexec/droid-hybris/system/etc/init/permissions_fix.rc
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/ld.config.28.txt
/usr/libexec/droid-hybris/system/lib/hw/sensors.sdm660.so
/usr/lib64/qt5/plugins/feedback/droid-vibrator-device.ini
/usr/lib/systemd/system/ldconfig-mount.service
/usr/lib/systemd/system/local-fs.target.wants/ldconfig-mount.service
/usr/share/csd/settings.d/hw-settings.ini
/usr/share/ssu/features.d/adaptation-community.ini
